# systemd/claude-tmux.service
[Unit]
Description=Claude Session Manager (Margatroid)
After=network.target

[Service]
Type=simple
Environment=PATH=%h/.cargo/bin:%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/bin/orchestrator-boot
ExecStop=/usr/bin/tmux kill-session -t claude
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
