# systemd/claude-daemon.service
[Unit]
Description=Claude Remote Control Daemon (Margatroid)
After=network.target claude-tmux.service

[Service]
Type=simple
Environment=PATH=%h/.cargo/bin:%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=-%h/.margatroid/install.sh
ExecStart=%h/bin/orchestrator-daemon
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
