# systemd/margatroid-tmux.service
[Unit]
Description=Margatroid Session Manager
After=network.target

[Service]
Type=simple
Environment=PATH=%h/.margatroid/bin:%h/.cargo/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.margatroid/bin/margatroid-boot
ExecStop=/usr/bin/tmux kill-session -t margatroid
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
