# systemd/margatroid-daemon.service
[Unit]
Description=Margatroid Remote Control Daemon
After=network.target margatroid-tmux.service

[Service]
Type=simple
Environment=PATH=%h/.margatroid/bin:%h/.cargo/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=-%h/.margatroid/install.sh
ExecStart=%h/.margatroid/bin/margatroid-daemon
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
