# systemd/margatroid-update.service
[Unit]
Description=Margatroid Auto-Update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PATH=%h/.margatroid/bin:%h/.cargo/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.margatroid/repo/install.sh
TimeoutStartSec=600
